# deploy/base-radar.service
[Unit]
Description=Base Smart Money Radar
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=robinhood-radar
Group=robinhood-radar
WorkingDirectory=/opt/robinhood-radar
Environment=NODE_ENV=production
Environment=BASE_HOST=127.0.0.1
Environment=BASE_PORT=18119
Environment=BASE_DATA_FILE=/var/lib/robinhood-radar/base.sqlite
Environment=BASE_MIN_ENTRY_USD=500
Environment=BASE_SIGNIFICANT_PROFIT_RATE=0.002
Environment=BASE_SMART_BASE_MULTIPLE=5
Environment=BASE_STRICT_MULTIPLE=10
Environment=BASE_REPEAT_MIN_HITS=2
Environment=BASE_STRONG_HOLDER_RANK=30
Environment=BASE_RELATED_CLUSTER_PENALTY=0.9
Environment=BASE_LOW_FREQUENCY_REASON_THRESHOLD=0.8
Environment=BASE_HOLDER_CANDIDATE_LIMIT=100
Environment=BASE_HOLDER_FETCH_LIMIT=150
Environment=BASE_HOLDER_PROFIT_CONCURRENCY=6
Environment=BASE_REQUEST_TIMEOUT_MS=20000
Environment=BASE_SCAN_CONCURRENCY=1
Environment=BASE_RPC_MAX_RETRIES=6
Environment=BASE_RPC_RETRY_DELAY_MS=500
Environment=BASE_RPC_MAX_RETRY_DELAY_MS=15000
Environment=BASE_RPC_BATCH_SIZE=10
Environment=BASE_RPC_BATCH_DELAY_MS=1000
Environment=BASE_MONITOR_POLL_INTERVAL_MS=500
Environment=BASE_MONITOR_DEGRADED_POLL_INTERVAL_MS=1000
Environment=BASE_MONITOR_MAX_BLOCK_SPAN=500
Environment=BASE_MONITOR_WALLET_TOPIC_CHUNK_SIZE=100
Environment=BASE_MONITOR_LOG_CONCURRENCY=2
Environment=BASE_MONITOR_RECOVERY_SUCCESSES=20
Environment=BASE_MONITOR_FAST_LIVE_BLOCK_SPAN=50
Environment=BASE_MONITOR_FAST_GAP_BLOCK_SPAN=100
Environment=BASE_MONITOR_FAST_GAP_POLL_INTERVAL_MS=5000
Environment=BASE_MONITOR_DEEP_POLL_INTERVAL_MS=500
Environment=BASE_MONITOR_DEEP_DEGRADED_POLL_INTERVAL_MS=1500
Environment=BASE_MONITOR_DEEP_LIVE_BLOCK_SPAN=20
Environment=BASE_MONITOR_DEEP_GAP_BLOCK_SPAN=20
Environment=BASE_MONITOR_DEEP_GAP_POLL_INTERVAL_MS=5000
Environment=BASE_MONITOR_TOKEN_METADATA_BUDGET_MS=1500
Environment=BASE_MARKET_DEBOT_BUDGET_MS=750
Environment=BASE_MARKET_REQUEST_TIMEOUT_MS=5000
ExecStart=/usr/local/bin/node /opt/robinhood-radar/base-server.mjs
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/robinhood-radar
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
