# deploy/robinhood-radar.service
[Unit]
Description=Robinhood Smart Money Radar
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=robinhood-radar
Group=robinhood-radar
WorkingDirectory=/opt/robinhood-radar
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=18118
Environment=ROBINHOOD_PUBLIC_DIR=/opt/robinhood-radar/public
Environment=ROBINHOOD_DATA_FILE=/var/lib/robinhood-radar/robinhood.sqlite
Environment=ROBINHOOD_MIN_ENTRY_USD=500
Environment=ROBINHOOD_SIGNIFICANT_PROFIT_RATE=0.002
Environment=ROBINHOOD_SMART_BASE_MULTIPLE=5
Environment=ROBINHOOD_STRICT_MULTIPLE=10
Environment=ROBINHOOD_REPEAT_MIN_HITS=2
Environment=ROBINHOOD_STRONG_HOLDER_RANK=30
Environment=ROBINHOOD_RELATED_CLUSTER_PENALTY=0.9
Environment=ROBINHOOD_LOW_FREQUENCY_REASON_THRESHOLD=0.8
Environment=ROBINHOOD_HOLDER_CANDIDATE_LIMIT=100
Environment=ROBINHOOD_HOLDER_FETCH_LIMIT=150
Environment=ROBINHOOD_HOLDER_PROFIT_CONCURRENCY=6
Environment=ROBINHOOD_REQUEST_TIMEOUT_MS=20000
Environment=ROBINHOOD_AUTO_SCAN_LIMIT=10
Environment=ROBINHOOD_MAX_SWAPS=1000
Environment=ROBINHOOD_SCAN_CONCURRENCY=1
Environment=ROBINHOOD_RPC_MAX_RETRIES=6
Environment=ROBINHOOD_RPC_RETRY_DELAY_MS=500
Environment=ROBINHOOD_RPC_MAX_RETRY_DELAY_MS=15000
Environment=ROBINHOOD_RPC_BATCH_SIZE=10
Environment=ROBINHOOD_RPC_BATCH_DELAY_MS=1000
Environment=ROBINHOOD_MONITOR_POLL_INTERVAL_MS=500
Environment=ROBINHOOD_MONITOR_DEGRADED_POLL_INTERVAL_MS=1000
Environment=ROBINHOOD_MONITOR_MAX_BLOCK_SPAN=500
Environment=ROBINHOOD_MONITOR_WALLET_TOPIC_CHUNK_SIZE=100
Environment=ROBINHOOD_MONITOR_LOG_CONCURRENCY=2
Environment=ROBINHOOD_MONITOR_RECOVERY_SUCCESSES=20
Environment=ROBINHOOD_MONITOR_DEEP_POLL_INTERVAL_MS=500
Environment=ROBINHOOD_MONITOR_DEEP_DEGRADED_POLL_INTERVAL_MS=1500
Environment=ROBINHOOD_MONITOR_DEEP_LIVE_BLOCK_SPAN=20
Environment=ROBINHOOD_MONITOR_DEEP_GAP_BLOCK_SPAN=20
Environment=ROBINHOOD_MONITOR_DEEP_GAP_POLL_INTERVAL_MS=5000
Environment=ROBINHOOD_MONITOR_TOKEN_METADATA_BUDGET_MS=1500
ExecStart=/usr/local/bin/node /opt/robinhood-radar/robinhood-server.mjs
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/robinhood-radar
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
