# deploy/solana-radar.service
[Unit]
Description=Solana Smart Money Radar
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=robinhood-radar
Group=robinhood-radar
WorkingDirectory=/opt/robinhood-radar
Environment=NODE_ENV=production
Environment=SOLANA_HOST=127.0.0.1
Environment=SOLANA_PORT=18120
Environment=SOLANA_DATA_FILE=/var/lib/robinhood-radar/solana.sqlite
Environment=SOLANA_MIN_ENTRY_USD=500
Environment=SOLANA_HOLDER_CANDIDATE_LIMIT=100
Environment=SOLANA_HOLDER_FETCH_LIMIT=150
Environment=SOLANA_HOLDER_PROFIT_CONCURRENCY=6
Environment=SOLANA_HOLDER_MAX_ACCOUNTS=50000
Environment=SOLANA_RPC_MAX_RESPONSE_BYTES=25165824
Environment=SOLANA_REQUEST_TIMEOUT_MS=20000
Environment=SOLANA_RPC_MAX_RETRIES=3
Environment=SOLANA_RPC_RETRY_DELAY_MS=500
Environment=SOLANA_RPC_MAX_RETRY_DELAY_MS=8000
Environment=SOLANA_SCAN_CONCURRENCY=1
Environment=SOLANA_WEBHOOK_BODY_LIMIT=4194304
EnvironmentFile=-/etc/robinhood-radar/solana.env
ExecStart=/usr/local/bin/node /opt/robinhood-radar/solana-server.mjs
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/robinhood-radar
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
